# rootfs/etc/services.d/react/run
#!/usr/bin/with-contenv bash

cd /app/react || exit
su dogen
yarn build
yarn run prod