# rootfs/etc/services.d/express/run
#!/usr/bin/with-contenv bash

su dogen
cd /app/express && yarn start start